# mesh/etc/rc.local
# Put your custom commands here that should be executed once
# the system init finished. By default this file does nothing.

/usr/bin/nexfi-led.sh > /dev/null &

# start l2cm_libuci process for topology
sleep 5
/usr/bin/l2cm_libuci > /dev/null &
